# miccpipe/systemd/micc-pipeline-runner.timer
[Unit]
Description=Run the MICC pipeline every 15 seconds

[Timer]
OnActiveSec=0
OnUnitActiveSec=15

[Install]
WantedBy=timers.target
